# api-service-main/install/websockets.conf
[program:RunPusherWebsocketListener]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan websockets:serve
autostart=true
autorestart=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/websocket-conf.log
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true
priority=99